# files/fancontroller
#!/bin/sh /etc/rc.common

START=99

start() {
  . /lib/functions.sh
  rmmod pwm_fan
}
